# deploy/nomad-eye-backend.service
[Unit]
Description=Nomad Eye Backend
After=network.target

[Service]
User=arduino
WorkingDirectory=/opt/nomad-eye/backend
Environment=PATH=/opt/nomad-eye/backend/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/nomad-eye/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8080 --loop uvloop
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
